# scripts/backup/systemd/alt-backup-offsite.timer
[Unit]
Description=Alt Platform Daily Offsite Sync Timer
Documentation=file:///home/koko/Documents/dev/Alt/docs/runbooks/backup-restore.md

[Timer]
# Run daily at 5:00 AM (after daily backup completes)
OnCalendar=*-*-* 05:00:00
# Randomize start within 30 minutes
RandomizedDelaySec=1800
# Persist timer across reboots
Persistent=true
# Accuracy within 1 minute
AccuracySec=1min

[Install]
WantedBy=timers.target
